# Software/Install/autostart/AstraInaHmi.desktop
[Desktop Entry]
Type=Application
Name=AstraIna
Exec=/opt/AstraDIY/AstraInaHmi.py

