# Software/Install/autostart/AstraInaHmi.desktop
[Desktop Entry]
Type=Application
Name=AstraIna
Exec=/opt/AstrAlim/AstraInaHmi.py

