# Software/Install/autostart/AstraGpsHmi.desktop
[Desktop Entry]
Type=Application
Name=AstraGps
Exec=/opt/AstraDIY/AstraGpsHmi.py

