# Software/Install/autostart/AstraDIYHmi.desktop
[Desktop Entry]
Type=Application
Name=AstraDIY
Exec=/opt/AstraDIY/AstraDIYHmi.py

