# Software/Install/autostart/AstraSuiteHmi.desktop
[Desktop Entry]
Type=Application
Name=AstraSuite
Exec=/opt/AstrAlim/AstraSuiteHmi.py

