# Software/Install/autostart/AstraGpioHmi.desktop
[Desktop Entry]
Type=Application
Name=AstraGpio
Exec=/opt/AstraDIY/AstraGpioHmi.py

