# Software/Install/autostart/AstraPwmHmi.desktop
[Desktop Entry]
Type=Application
Name=AstraPwm
Exec=/opt/AstrAlim/AstraPwmHmi.py

